# scripts/deploy/supervisor/super_uwsgi.conf
[program:uwsgi]
command = newrelic-admin run-program uwsgi --ini /henrique/scripts/deploy/uwsgi/henrique.ini
